# gguf-switchboard.service
[Unit]
Description=GGUF Switchboard — GPU-aware local GGUF model scheduler
After=network.target
Documentation=https://github.com/pradeepgudipati/gguf-switchboard

[Service]
Type=simple
User=gguf-switchboard
Group=gguf-switchboard
ExecStart=/usr/local/bin/gguf-switchboard /etc/gguf-switchboard/config.toml
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
Environment=RUST_LOG=info
EnvironmentFile=-/etc/gguf-switchboard/env

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/models
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
